# sample/rc.local
LXC_NETWORK=10.0.0

# ONOS
iptables -t nat -A PREROUTING -p tcp --dport 8181 -j DNAT --to $LXC_NETWORK.11:8181
iptables -A INPUT -m state --state ESTABLISHED,RELATED -j ACCEPT
